# script/amfsanfrontend.service
[Unit]
Description=AMFSAN Frontend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/app/.next/standalone
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXT_PUBLIC_API_BASE_URL=http://your-backend-api:8080

# 资源限制
LimitNOFILE=65536
MemoryMax=200M
CPUQuota=50%

[Install]
WantedBy=multi-user.target